# runit/offlineimap/run
#!/bin/sh

set -eu

if [ "${DEBUG:-}" ]; then
	set -- -u blinkenlights -o
else
	set -- -u basic
fi

exec 2>&1
exec chpst -e "$(dirname "$0")/env" /usr/bin/env offlineimap $@
